# src/common/supervisor/butterfly.ini
[program:butterfly]
command=/usr/bin/butterfly.server.py --host=127.0.0.1 --port=57575 --unsecure --login --motd=/etc/motd %(ENV_BUTTERFLY_OPTSOPTS)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
