# src/common/supervisor/filebrowser.ini
[program:filer]
command=filebrowser -r /root -a 127.0.0.1 -p 57576 -d /var/lib/filebrowser/filebrowser.db -b /file --noauth
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
