# src/common/supervisor/ttyd.ini
[program:ttyd]
command=/usr/local/bin/ttyd %(ENV_TTYD_OPTS)s -p 57575 -i lo -b /term/ /bin/bash
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
